# services/friday-core.service
[Unit]
Description=Friday Core API
After=network-online.target friday-vllm.service
Wants=network-online.target
Requires=friday-vllm.service

[Service]
Type=simple
WorkingDirectory=/home/artur/friday
Environment="PATH=/home/artur/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/artur/friday"

# Kill all child processes when stopping
KillMode=control-group
KillSignal=SIGKILL
TimeoutStopSec=10

# Ensure port is free and vLLM is ready before starting
ExecStartPre=/bin/bash -c 'fuser -k 8080/tcp 2>/dev/null || true; sleep 1'
ExecStartPre=/bin/bash -c 'for i in {1..30}; do curl -s http://localhost:8000/v1/models >/dev/null && exit 0; sleep 3; done; exit 1'

ExecStart=/home/artur/.local/bin/pipenv run uvicorn src.api.routes:app --host 0.0.0.0 --port 8080 --log-level info

# Allow enough time for health checks to complete
TimeoutStartSec=120

Restart=always
RestartSec=10
StandardOutput=append:/home/artur/friday/logs/friday-core.log
StandardError=append:/home/artur/friday/logs/friday-core.log

[Install]
WantedBy=default.target
