# services/friday-telegram.service
[Unit]
Description=Friday Telegram Bot
After=network-online.target friday-core.service
Wants=network-online.target
Requires=friday-core.service

[Service]
Type=simple
WorkingDirectory=/home/artur/friday
Environment="PATH=/home/artur/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/artur/friday"
ExecStart=/home/artur/.local/bin/pipenv run python -m src.telegram_bot
Restart=always
RestartSec=10
StandardOutput=append:/home/artur/friday/logs/friday-telegram.log
StandardError=append:/home/artur/friday/logs/friday-telegram.log

[Install]
WantedBy=default.target
