# services/homelab-monitor.service
[Unit]
Description=Friday AI Homelab Monitor
After=network.target friday.service

[Service]
Type=simple
User=artur
WorkingDirectory=/home/artur/friday
Environment="PATH=/home/artur/.local/share/virtualenvs/friday/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/artur/.local/share/virtualenvs/friday/bin/python -m src.awareness
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/artur/friday/logs/homelab_monitor.log
StandardError=append:/home/artur/friday/logs/homelab_monitor.log

[Install]
WantedBy=multi-user.target
